# com.stockdashboard.fetchdaily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockdashboard.fetchdaily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/r/Documents/Claude/stock-dashboard/run_fetch_daily.sh</string>
    </array>

    <!-- 평일(월~금) 16:05 KST 실행 (장 마감 35분 후) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>5</integer>
        </dict>
    </array>

    <!-- Mac이 잠자기 상태였다가 깨어나면 놓친 실행을 한 번 수행 -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/r/Documents/Claude/stock-dashboard/logs/fetch_daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/r/Documents/Claude/stock-dashboard/logs/fetch_daily_err.log</string>
</dict>
</plist>
